# libfranka-sim/k8s/desktop.conf
[program:vnc]
command=/usr/bin/vncserver -geometry 1280x800 -localhost no -fg :1
user=ubuntu
autorestart=true
priority=10

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc/ 6080 localhost:5901
autorestart=true
priority=20